# apps/worker/deploy/outside-ir35-jobs-worker.service
[Unit]
Description=Outside IR35 Jobs Aggregation Worker
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/outside-ir35-jobs/apps/worker
EnvironmentFile=/opt/outside-ir35-jobs/apps/worker/.env
# xvfb-run gives a headless X display for the local Playwright fallback. When
# scraping via Browserbase (remote browser) it's harmless. Run via Bun (Hetzner).
ExecStart=/usr/bin/xvfb-run --auto-servernum --server-args="-screen 0 1920x1080x24" /root/.bun/bin/bun run src/index.ts
Restart=always
RestartSec=5
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
